# operations/config/maingo.conf
[program:maingo]
command=/usr/bin/go run /home/ec2-user/main.go
user=ec2-user
stopasgroup = true
killasgroup = true
